# scripts/com.bucketsvision.snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bucketsvision.snapshot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/stlee/Desktop/bucketsvision/scripts/auto_snapshot.sh</string>
    </array>

    <!-- 매일 오후 6시 (KST) 실행 = 오전 4시 (ET) -->
    <!-- NBA 경기는 보통 오전 1-2시 (ET)에 종료 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/stlee/Desktop/bucketsvision/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/stlee/Desktop/bucketsvision/logs/launchd_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
